# server_configs/init/face_to_face.conf
# upstart service file at /etc/init/C5.conf
description "Meteor.js (NodeJS) application"
author "Don Undeen <donundeen@gmail.com>"

# When to start the service
start on started mongodb and runlevel [2345]

# When to stop the service
stop on shutdown

# Automatically restart process if crashed
respawn
respawn limit 10 5

# we don't use buil-in log because we use a script below
# console log

# drop root proviliges and switch to mymetorapp user
setuid donundeen
setgid donundeen

script
    export PATH=/home/donundeen/.nvm/v0.10.40/bin/:/opt/local/bin:/opt/local/sbin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
    export NODE_PATH=/home/donundeen/.nvm/v0.10.40/bin/node:/home/donundeen/.nvm/v0.10.40/lib/node_modules:/usr/share/javascript
    # set to home directory of the user Meteor will be running as
    export PWD=/home/donundeen
    export HOME=/home/donundeen
    # leave as 127.0.0.1 for security
    export BIND_IP=127.0.0.1
    # the port nginx is proxying requests to
    export PORT=9010
    # this allows Meteor to figure out correct IP address of visitors
    export HTTP_FORWARDED_COUNT=1
    # MongoDB connection string using todos as database name
    export MONGO_URL=mongodb://localhost:27017/meteor
    # The domain name as configured previously as server_name in nginx
    export ROOT_URL=http://208.43.234.91
    # optional JSON config - the contents of file specified by passing "--settings" parameter to meteor command in development mode
    export METEOR_SETTINGS='{ "somesetting": "someval", "public": { "othersetting": "anothervalue" } }'
    # this is optional: http://docs.meteor.com/#email
    # commented out will default to no email being sent
    # you must register with MailGun to have a username and password there
    # export MAIL_URL=smtp://postmaster@mymetorapp.net:password123@smtp.mailgun.org
    # alternatively install "apt-get install default-mta" and uncomment:
    # export MAIL_URL=smtp://localhost
    exec  node /home/donundeen/face_to_face_deploy/bundle/main.js >> /home/donundeen/face_to_face_meteor.log
end script